# launchd/com.anyshell.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anyshell.web</string>

    <key>ProgramArguments</key>
    <array>
        <string>HOMEDIR/.local/bin/ttyd-wrapper</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ANYSHELL_SESSION</key>
        <string>main</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:HOMEDIR/.local/bin</string>
        <key>HOME</key>
        <string>HOMEDIR</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>HOMEDIR/.local/share/anyshell/ttyd.log</string>

    <key>StandardErrorPath</key>
    <string>HOMEDIR/.local/share/anyshell/ttyd.error.log</string>

    <key>WorkingDirectory</key>
    <string>HOMEDIR</string>

    <key>Umask</key>
    <integer>77</integer>
</dict>
</plist>
